# crontab
# BBC Noticias bot — runs daily at 08:00 CET/CEST
TZ=Europe/Berlin
0 8 * * * root cd /app && uv run python -m src.bbc_noticias.bot >> /app/data/cron.log 2>&1