# crontab
# BBC Noticias bot — runs daily at 08:00 CET/CEST
CRON_TZ=Europe/Berlin
0 8 * * * root cd /app && uv run python -m src.bbc_noticias.bot